# addon/rootfs/etc/s6-overlay/s6-rc.d/captive-portal-guest/run
#!/command/with-contenv bashio
# SPDX-FileCopyrightText: 2026 Andrew Grimberg
# SPDX-License-Identifier: Apache-2.0
# shellcheck shell=bash
# s6-overlay run script for the guest captive portal listener.

# Omada controller settings and guest_external_url are now managed via
# the web UI and stored in the database.  No environment variable
# exports are needed for those.

declare debug_guest_portal
debug_guest_portal=$(bashio::config 'debug_guest_portal' 'false')
export CP_DEBUG_GUEST_PORTAL="${debug_guest_portal}"

bashio::log.info "Starting guest captive portal listener on port 8099..."

exec python -m uvicorn \
    captive_portal.guest_app:create_guest_app \
    --factory \
    --host 0.0.0.0 \
    --port 8099
