# addon/rootfs/etc/s6-overlay/s6-rc.d/captive-portal-guest/run
#!/command/with-contenv bashio
# SPDX-FileCopyrightText: 2026 Andrew Grimberg
# SPDX-License-Identifier: Apache-2.0
# shellcheck shell=bash
# s6-overlay run script for the guest captive portal listener.

declare guest_external_url
guest_external_url=$(bashio::config 'guest_external_url' '')
export CP_GUEST_EXTERNAL_URL="${guest_external_url}"

# Export Omada controller environment variables for AppSettings
declare omada_controller_url
omada_controller_url=$(bashio::config 'omada_controller_url' '')
export CP_OMADA_CONTROLLER_URL="${omada_controller_url}"

declare omada_username
omada_username=$(bashio::config 'omada_username' '')
export CP_OMADA_USERNAME="${omada_username}"

declare omada_password
omada_password=$(bashio::config 'omada_password' '')
export CP_OMADA_PASSWORD="${omada_password}"

declare omada_site_name
omada_site_name=$(bashio::config 'omada_site_name' '')
export CP_OMADA_SITE_NAME="${omada_site_name}"

declare omada_controller_id
omada_controller_id=$(bashio::config 'omada_controller_id' '')
export CP_OMADA_CONTROLLER_ID="${omada_controller_id}"

declare omada_verify_ssl
omada_verify_ssl=$(bashio::config 'omada_verify_ssl' '')
export CP_OMADA_VERIFY_SSL="${omada_verify_ssl}"

declare debug_guest_portal
debug_guest_portal=$(bashio::config 'debug_guest_portal' 'false')
export CP_DEBUG_GUEST_PORTAL="${debug_guest_portal}"

bashio::log.info "Starting guest captive portal listener on port 8099..."
if [ -n "${guest_external_url}" ]; then
    bashio::log.info "Guest external URL: ${guest_external_url}"
fi

exec python -m uvicorn \
    captive_portal.guest_app:create_guest_app \
    --factory \
    --host 0.0.0.0 \
    --port 8099
