# addon/rootfs/etc/s6-overlay/s6-rc.d/captive-portal/run
#!/command/with-contenv bashio
# SPDX-FileCopyrightText: 2025 Andrew Grimberg
# SPDX-License-Identifier: Apache-2.0
# shellcheck shell=bash
# s6-overlay run script for the captive portal service.

declare log_level
log_level=$(bashio::config 'log_level' 'info')
bashio::log.info "Starting captive portal (log_level=${log_level})..."

# Retrieve the ingress entry path so the app can generate correct URLs
# when served behind the Home Assistant ingress proxy.
declare ingress_entry
ingress_entry=$(bashio::addon.ingress_entry)
bashio::log.info "Ingress entry configured for captive portal ingress."

# Export Omada controller environment variables for AppSettings
declare omada_controller_url
omada_controller_url=$(bashio::config 'omada_controller_url' '')
export CP_OMADA_CONTROLLER_URL="${omada_controller_url}"

declare omada_username
omada_username=$(bashio::config 'omada_username' '')
export CP_OMADA_USERNAME="${omada_username}"

declare omada_password
omada_password=$(bashio::config 'omada_password' '')
export CP_OMADA_PASSWORD="${omada_password}"

declare omada_site_name
omada_site_name=$(bashio::config 'omada_site_name' '')
export CP_OMADA_SITE_NAME="${omada_site_name}"

declare omada_controller_id
omada_controller_id=$(bashio::config 'omada_controller_id' '')
export CP_OMADA_CONTROLLER_ID="${omada_controller_id}"

declare omada_verify_ssl
omada_verify_ssl=$(bashio::config 'omada_verify_ssl' '')
export CP_OMADA_VERIFY_SSL="${omada_verify_ssl}"

declare debug_guest_portal
debug_guest_portal=$(bashio::config 'debug_guest_portal' 'false')
export CP_DEBUG_GUEST_PORTAL="${debug_guest_portal}"

exec python -m uvicorn \
    captive_portal.app:create_app \
    --factory \
    --host 0.0.0.0 \
    --port 8080 \
    --root-path "${ingress_entry}"
